# config/systemd/user/emacs.socket
[Socket]
ListenStream=%t/emacs/server
SocketMode=0700
DirectoryMode=0700

[Install]
WantedBy=sockets.target
