# config/systemd/user/mpd.socket
[Socket]
ListenStream=%t/mpd/socket
SocketMode=0700
DirectoryMode=0700

[Install]
WantedBy=sockets.target
