# systemd-user/bluetooth-applet.service
[Unit]
Description=Bluetooth Manager Applet
Requisite=bluetooth.target

[Service]
ExecStart=@/usr/bin/bluetooth-applet bluetooth-applet
