# systemd-user/dbus.socket
[Unit]
Description=D-Bus Session (socket)

[Socket]
ListenStream=%t/dbus/user_bus_socket

[Install]
WantedBy=default.target