# systemd-user/dropboxd.service
[Unit]
Description=run dropbox

[Service]
ExecStart=/usr/bin/dropboxd

[Install]
WantedBy=default.target
