# systemd-user/envoy@.socket
[Socket]
ListenStream=@/vodik/envoy

[Install]
WantedBy=default.target
