# systemd-user/i3.target
[Unit]
Description=i3 Desktop Session

Requires=desktop-basic.target
Requires=i3.service

[Install]
WantedBy=default.target
