# systemd-user/mpd.service
[Unit]
Description=Music Player Daemon
Requires=pulse.service
Requires=mpd.socket

[Service]
ExecStart=/usr/bin/mpd --stdout --no-daemon
