# systemd-user/nm-applet.service
[Unit]
Description=NetworkManager Applet

[Service]
BusName=org.gnome.network_manager_applet
ExecStart=@/usr/bin/nm-applet nm-applet
