# systemd-user/udiskie.service
[Unit]
Description=handle automounting
Requires=dbus.service

[Service]
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target