# systemd-user/xbindkeys.service
[Unit]
Description=xrdb
#Requires=xorg.target
#After=xorg.target

[Service]
ExecStart=/usr/bin/xbindkeys --nodaemon
